# .config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/$(whoami)/.config/systemd/user/check-battery.sh

