# .config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/$USER/.config/systemd/user/check-battery.sh

