# .config/systemd/user/check-battery.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min
Unit=check-battery.service

[Install]
WantedBy=timers.target

