# .config/systemd/user/mailsync.timer
[Unit]
Description=Timer to sync emails

[Timer]
OnActiveSec=10min
OnUnitActiveSec=10min
Unit=mailsync.service

[Install]
WantedBy=timers.target

