# .config/systemd/user/mailsync.service
[Unit]
Description=mailsync

[Service]
ExecStart=/usr/bin/mailsync

