# .config/systemd/user/theme-switcher.service
[Unit]
Description=Theme switcher

[Service]
ExecStart=/home/$(whoami)/.config/systemd/user/theme-switcher.sh


